# ressources/pkt_chksvc.service
[Unit]
Description=PKT Check services
After=network.target

[Service]
ExecStart=/bin/bash TARGETPATH/pkt_chksvc.sh
Restart=always

[Install]
WantedBy=multi-user.target